# packages/auto-knight/src/usr/lib/systemd/user/auto-knight.service
[Unit]
Description=Automatically toggle Plasma dark mode
PartOf=plasma-plasmashell.service
After=plasma-plasmashell.service

[Service]
Type=simple
ExecStart=/usr/bin/auto-knight
Restart=on-failure
RestartSec=10s

[Install]
UpheldBy=plasma-plasmashell.service
