# packages/evanc577-update-services/src/usr/lib/systemd/system/flatpak-system-update.service.d/z-override.conf
[Service]
ExecCondition=
TimeoutStartSec=infinity
ExecStartPre=/usr/bin/wait-for-network
